# crontab
# min hour day month weekday command
0 1 * * * su node -c 'node /opt/membership-system/built/tools/mailchimp/sync.js'
0 1 * * * su node -c 'node /opt/membership-system/built/tools/start-gifts.js'
